# d369.service
[Unit]
Description=d369 — يقرأ الكتاب والكون
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/d369
ExecStart=/usr/bin/python3 /root/d369/bot.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
